# server/supervisord.conf
[supervisord]
nodaemon=true

[program:django]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/bin/bash ./start.sh
user=polyledger

[program:celery_worker]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=celery -A polyledger worker -l info
user=polyledger

[program:celery_beat]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=celery -A polyledger beat -l info -s celerybeat-schedule.db
user=polyledger
